# devopsfetch.service
[Unit]
Description=DevopsFetch Service
After=network.target

[Service]
ExecStart=$SCRIPT_PATH $SCRIPT_OPTIONS
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target